# templates/multizone-audio-control.service.template
[Unit]
Description=Multizone Audio Controller
After=network.target

[Service]
ExecStart={{python_venv}} {{config_root}}/controller/multizone-control.py
#User=mopidy
#Group=mopidy
Restart=on-failure

[Install]
WantedBy=multi-user.target
